# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/app/venv/bin/python app.py

[program:prometheus]
command=/app/venv/bin/python prometheus/app.py